# supervisor/nginx/templates/nginx.ini.j2
[program:nginx]
command = /usr/sbin/nginx -c /etc/nginx/nginx.conf
user = nginx
stdout_events_enabled = true
stderr_events_enabled = true
autostart = true
autorestart = true
stopasgroup = true
