# supervisor/postfix/templates/postfix.ini.j2
[program:postfix]
command = {{ supervisor_exec_shell_dir }}/postfix.sh
stdout_events_enabled = true
stderr_events_enabled = true
autostart = true
autorestart = true
stopasgroup = true
